# deploy/_supervisor.conf
[program:]
process_name=%(program_name)s
command=
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=
numprocs=8
redirect_stderr=true
stdout_logfile=
stopwaitsecs=3600